# deployments/service01/templates.bak.20250822_081607/service01.service.j2
[Unit]
Description={{ app_name }} - Bun Application
Documentation=https://bun.sh/
After=network.target
Wants=network.target

[Service]
Type=simple
User={{ app_user }}
Group={{ app_user }}
WorkingDirectory={{ app_dir }}
ExecStart=/usr/local/bin/bun index.js
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier={{ app_service_name }}

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{ app_dir }}
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Environment
Environment=NODE_ENV={{ production_env.NODE_ENV }}
Environment=PORT={{ production_env.PORT }}
Environment=HOST={{ production_env.HOST }}
EnvironmentFile=-{{ app_dir }}/.env

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
